# units/arwain.socket
[Socket]
ListenFIFO=/tmp/arwain.stdin
Service=arwain.service

[Install]
WantedBy=sockets.target
